# dd/rootfs/etc/services.d/dd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the dd service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
bashio::log.info "Starting dd..."

# Retrieve config/service values
DD_HOST="$(bashio::config 'host')"
MQTT_PREFIX="$(bashio::config 'mqtt_prefix')"
DEBUG=$(bashio::config 'debug')

# Check if user provided custom MQTT settings or use HA MQTT service
if bashio::config.exists 'mqtt.broker'; then
    bashio::log.info "Using custom MQTT broker configuration."
    MQTT_HOST="$(bashio::config 'mqtt.broker')"
    MQTT_PORT="$(bashio::config 'mqtt.port')"
    MQTT_USER="$(bashio::config 'mqtt.username')"
    MQTT_PASSWORD="$(bashio::config 'mqtt.password')"
else
    bashio::log.info "Using Home Assistant MQTT service."

    # Wait for MQTT service to become available
    if ! bashio::services.available 'mqtt'; then
        bashio::log.warning "Home Assistant MQTT service not available; waiting..."
        bashio::services.wait_for 'mqtt'
        bashio::log.info "MQTT service is now available."
    fi

    MQTT_HOST=$(bashio::services 'mqtt' 'host')
    MQTT_PORT=$(bashio::services 'mqtt' 'port')
    MQTT_USER=$(bashio::services 'mqtt' 'username')
    MQTT_PASSWORD=$(bashio::services 'mqtt' 'password')
fi

# Conditionally set the debug flag
if [ "${DEBUG}" = "true" ]; then
    DEBUG_FLAG="-debug"
    bashio::log.info "Debug mode is enabled."
else
    DEBUG_FLAG=""
fi

# Run your dd/haus command, appending the debug flag if set
exec /usr/bin/dd/haus \
    -host "${DD_HOST}" \
    -mqtt "${MQTT_HOST}" \
    -mqttPort "${MQTT_PORT}" \
    -mqttUser "${MQTT_USER}" \
    -mqttPassword "${MQTT_PASSWORD}" \
    -mqttPrefix "${MQTT_PREFIX}" \
    -credentials /config/dd-credentials.json \
    ${DEBUG_FLAG}
