# dd/rootfs/etc/services.d/dd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the dd service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
bashio::log.info "Starting dd..."

# Wait for MQTT service to become available rather than exiting
if ! bashio::services.available 'mqtt'; then
    bashio::log.warning "Home Assistant MQTT service not available; waiting..."
    bashio::services.wait_for 'mqtt'
    bashio::log.info "MQTT service is now available."
fi

bashio::log.info "Using Home Assistant MQTT service."

# Retrieve config/service values
DD_HOST="$(bashio::config 'host')"
MQTT_HOST=$(bashio::services 'mqtt' 'host')
MQTT_PORT=$(bashio::services 'mqtt' 'port')
MQTT_USER=$(bashio::services 'mqtt' 'username')
MQTT_PASSWORD=$(bashio::services 'mqtt' 'password')
DEBUG=$(bashio::services 'debug')    # or bashio::config 'debug' if your add-on config uses "debug" key

# Conditionally set the debug flag
if [ "${DEBUG}" = "true" ]; then
    DEBUG_FLAG="-debug"
    bashio::log.info "Debug mode is enabled."
else
    DEBUG_FLAG=""
fi

# Run your dd/haus command, appending the debug flag if set
exec /usr/bin/dd/haus \
    -host "${DD_HOST}" \
    -mqtt "${MQTT_HOST}" \
    -mqttPort "${MQTT_PORT}" \
    -mqttUser "${MQTT_USER}" \
    -mqttPassword "${MQTT_PASSWORD}" \
    -credentials /config/dd-credentials.json \
    ${DEBUG_FLAG}
