# dd/rootfs/etc/services.d/dd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the dd service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting dd..."

# Check if MQTT service is available
if bashio::services.available 'mqtt'; then
    bashio::log.info "Using Home Assistant MQTT service."

    # Log retrieved MQTT values for debugging
    DD_HOST="$(bashio::config 'host')"
    MQTT_HOST=$(bashio::services 'mqtt' 'host')
    MQTT_PORT=$(bashio::services 'mqtt' 'port')
    MQTT_USER=$(bashio::services 'mqtt' 'username')
    MQTT_PASSWORD=$(bashio::services 'mqtt' 'password')
    DEBUG=$(bashio::services 'debug')



    # Execute the service with the MQTT configuration
    exec /usr/bin/dd/haus \
        -host "${DD_HOST}" \
        -mqtt "${MQTT_HOST}" \
        -mqttPort "${MQTT_PORT}" \
        -mqttUser "${MQTT_USER}" \
        -mqttPassword "${MQTT_PASSWORD}" \
        -creds /config/creds.json \
        -debug "${DEBUG}"
else
    bashio::log.warning "MQTT configuration not found and Home Assistant MQTT service unavailable."
    bashio::log.error "Ensure the MQTT service is enabled in Home Assistant."
    exit 1
fi
